# Web/web/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=service apache2 start

[program:mysql]
command=service mysql start
